# indego.cron
# m     h   dom mon dow command
*/10    *   *   *   *   ~/venv/bin/python3 ~/app/psql_insert.py > ~/psql_insert.log 2>&1
